# deploy/docker/monitoring/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ---- munin-node (data collection daemon) ----
[program:munin-node]
command=/usr/sbin/munin-node --foreground
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ---- muninweb (HTTP face: static site + graph CGI) ----
# Binds 0.0.0.0 so nginx (a separate container) can reach it over the
# internal network; the port is exposed, never published to the host.
# Authentication lives in nginx, same as bare metal.
[program:muninweb]
command=/usr/local/lib/naust/muninweb -listen 0.0.0.0:4948
user=munin
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=15

# ---- munin (cron-style graph generator, run periodically via a loop) ----
# Munin has no foreground daemon mode; we wrap it in a simple loop that runs
# munin-cron every 5 minutes - matching the default cron schedule.
[program:munin]
command=/bin/bash -c 'while true; do sudo -H -u munin munin-cron; sleep 300; done'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
