# deploy/docker/dns/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ---- NSD (authoritative) ----
[program:nsd]
command=/usr/sbin/nsd -d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ---- unbound (recursive resolver) ----
[program:unbound]
command=/usr/sbin/unbound -d
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

# ---- Control socket server ----
# Accepts service lifecycle requests from the management container over a
# Unix socket and dispatches to /etc/naust/handlers/<service> <action>.
# Returns synchronous OK/ERROR so management gets real error propagation.
[program:control-socket]
command=/bin/bash -c 'source /etc/naust.conf && exec /usr/local/sbin/control-socket-server "$STORAGE_ROOT/sockets/dns.sock" /etc/naust/handlers'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50
