# deploy/docker/management/supervisord.conf
[supervisord]
# Log to stdout so Docker captures all output via 'docker logs'.
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ---- helperd (privileged; root) ----
# Runs the fixed intent menu (daemon/internal/helper). Must be listening
# before managerd dials it - lower priority number starts it first.
[program:helperd]
command=/usr/local/lib/naust/helperd -socket /run/naust/helper.sock -socket-group naust
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

# ---- managerd (unprivileged control plane) ----
# Bound to 0.0.0.0 so nginx (in a separate container) can reach it over the
# Docker bridge. -submit-addr points at the mail container, not localhost.
[program:managerd]
command=/bin/bash -c 'source /etc/naust.conf && exec /usr/local/lib/naust/managerd -listen 0.0.0.0:10223 -submit-addr "${MAIL_HOST:-mail}:25"'
user=naust
environment=HOME="/var/lib/naust-managerd",XDG_CACHE_HOME="/var/lib/naust-managerd/cache"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
